# build/xkeys-server.service
[Unit]
Description=Run xkeys-server as background service
After=mosquitto.target

[Service]
ExecStart=/opt/xkeys-server/xkeys-server-%ARCH%.AppImage

Restart=always
User=nobody
Group=%NOGROUP%

Environment=PATH=/usr/bin:/usr/local/bin
Environment=NODE_ENV=production

WorkingDirectory=/opt/xkeys-server

[Install]
WantedBy=multi-user.target

