# build/macos/com.xkeys-server.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
    <dict>

        <key>Label</key>
        <string>com.__PRODUCT__.daemon</string>

        <key>RunAtLoad</key>
        <true/>

	<key>KeepAlive</key>
	<true/>

        <key>StandardErrorPath</key>
        <string>/var/log/__PRODUCT__.log</string>

        <key>StandardOutPath</key>
        <string>/var/log/__PRODUCT__.log</string>

        <key>WorkingDirectory</key>
        <string>/Library/__PRODUCT__</string>

        <key>ProgramArguments</key>
        <array>
            <string>/Library/__PRODUCT__/__VERSION__/__PRODUCT__</string>
        </array>

    </dict>
</plist>

