# build/macos/com.xkeys-server.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
    <dict>

        <key>Label</key>
        <string>com.xkeys-server.daemon</string>

        <key>RunAtLoad</key>
        <true/>

	<key>KeepAlive</key>
	<true/>

        <key>StandardErrorPath</key>
        <string>/var/log/xkeys-server.log</string>

        <key>StandardOutPath</key>
        <string>/var/log/xkeys-server.log</string>

        <key>EnvironmentVariables</key>
        <dict>
            <key>PATH</key>
            <string><![CDATA[/usr/local/bin:/usr/local/sbin:/usr/bin:/bin:/usr/sbin:/sbin]]></string>
        </dict>

        <key>WorkingDirectory</key>
        <string>/Library/xkeys-server/__VERSION__</string>

        <key>ProgramArguments</key>
        <array>
            <string>/Library/xkeys-server/__VERSION__/xkeys-server</string>
        </array>

    </dict>
</plist>

